# debian/orchestrator/DEBIAN/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    if systemctl is-active --quiet mmorpg-orchestrator.service; then
        systemctl stop mmorpg-orchestrator.service || true
    fi
fi
if [ "$1" = "remove" ]; then
    systemctl disable mmorpg-orchestrator.service || true
fi
exit 0
