# debian/redis/DEBIAN/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    if systemctl is-active --quiet mmorpg-redis.service; then
        systemctl stop mmorpg-redis.service || true
    fi
fi
if [ "$1" = "remove" ]; then
    systemctl disable mmorpg-redis.service || true
fi
exit 0
